# gs232-bridge@.service
[Unit]
Description=GS-232 telescope bridge
Documentation=https://github.com/Josh3-14159/rp2040-gpio-fs
# Start after the FUSE filesystem is mounted (udev+systemd template unit)
After=network.target
# Uncomment once rp2040-gpio-fs is installed and confirmed working:
# After=network.target rp2040fs@%i.service
# BindsTo=rp2040fs@%i.service

[Service]
Type=simple
User=%i
RuntimeDirectory=gs232-bridge
RuntimeDirectoryMode=0755
ExecStart=/usr/bin/python3 /opt/gs232_bridge/main.py --config /opt/gs232_bridge/config.ini
# Create /dev/ttyGS232 symlink as root so SDRAngel can find it in /dev/tty*
ExecStartPost=/usr/local/bin/gs232-mklink
# Remove /dev/ttyGS232 on stop
ExecStopPost=/bin/rm -f /dev/ttyGS232

# SIGHUP reloads config without restart
ExecReload=/bin/kill -HUP $MAINPID

# Clean restart on unexpected crash (not on clean exit or SIGTERM)
Restart=on-failure
RestartSec=5

# Give the control loop time to stop motors cleanly before SIGKILL
TimeoutStopSec=10

# Logging goes to journald automatically
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gs232-bridge

[Install]
WantedBy=multi-user.target
